# ansible/roles/portctl/templates/portctl-web-local.service.j2
[Unit]
Description=Port Forward Web UI (localhost)
After=portctl.service
Requires=portctl.service

[Service]
Type=simple
WorkingDirectory={{ portctl_web_root }}
ExecStart=/usr/bin/php -S {{ portctl_web_local_bind }}:{{ portctl_web_port }} -t {{ portctl_web_root }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
